# output/conf/supervisord.conf
[unix_http_server]
file=/home/work/bfserver/run/supervisor.sock
chmod=0700
chown=work:work

[supervisord]
logfile=/data/logs/bfserver_supervisord.log
logfile_maxbytes=50MB  
logfile_backups=10     
loglevel=info          
pidfile=/home/work/bfserver/run/supervisord.pid 
nodaemon=false
minfds=1024  
minprocs=200 

[supervisorctl]
serverurl=unix:///home/work/bfserver/run/supervisor.sock 

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:bfserver]
command=/home/work/bfserver/bin/bfserver
process_name=bfserver
numprocs=1
directory=/home/work/bfserver/
stdout_logfile=/data/logs/bfserver/bfserver_stdout.log
stderr_logfile=/data/logs/bfserver/bfserver_stderr.log
stopsignal=TERM
stopwaitsecs=600
